# vocalyx-postgres.service
# ============================================================================
# VOCALYX - PostgreSQL Service pour systemd
# ============================================================================
# Service systemd pour gérer le conteneur PostgreSQL avec Podman
# 
# Installation:
#   sudo cp vocalyx-postgres.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable vocalyx-postgres.service
#   sudo systemctl start vocalyx-postgres.service
# ============================================================================

[Unit]
Description=PostgreSQL Database for Vocalyx
Documentation=https://www.postgresql.org/
After=network-online.target
Wants=network-online.target
Requires=vocalyx-postgres-data.service

[Service]
Type=notify
NotifyAccess=all
TimeoutStartSec=0
Restart=always
RestartSec=10

# Créer le conteneur s'il n'existe pas, sinon le démarrer
ExecStartPre=-/usr/bin/podman rm -f vocalyx-postgres
ExecStart=/usr/bin/podman run \
    --name vocalyx-postgres \
    --network vocalyx-network \
    --volume vocalyx-postgres-data:/var/lib/postgresql/data:Z \
    --env POSTGRES_DB=vocalyx_db \
    --env POSTGRES_USER=vocalyx \
    --env POSTGRES_PASSWORD=vocalyx_secret \
    --publish 5432:5432 \
    --memory 1g \
    --health-cmd "pg_isready -U vocalyx -d vocalyx_db" \
    --health-interval 10s \
    --health-timeout 5s \
    --health-retries 5 \
    --rm \
    postgres:15-alpine

ExecStop=/usr/bin/podman stop -t 10 vocalyx-postgres
ExecStopPost=-/usr/bin/podman rm -f vocalyx-postgres

[Install]
WantedBy=multi-user.target

