# vocalyx-api-02.container
# ============================================================================
# VOCALYX - API Container #2 pour Podman/systemd
# ============================================================================
# Fichier quadlet pour le conteneur vocalyx-api-02
# 
# Note: L'image doit être construite au préalable:
#   podman build -t vocalyx-api:latest -f ./vocalyx-api/Containerfile ./vocalyx-api
# 
# Installation:
#   sudo cp vocalyx-api-02.container /etc/containers/systemd/
#   sudo systemctl daemon-reload
#   sudo systemctl start vocalyx-api-02.service
# ============================================================================

[Unit]
Description=Vocalyx API Service #2
After=vocalyx-network.service vocalyx-postgres.service vocalyx-redis.service
Requires=vocalyx-network.service vocalyx-postgres.service vocalyx-redis.service

[Container]
Image=vocalyx-api:latest
ContainerName=vocalyx-api-02
RestartPolicy=unless-stopped
Network=vocalyx-network
Volume=%E/shared/uploads:/app/shared_uploads:Z
Volume=%E/shared/logs:/app/logs:Z
Memory=1G
Environment=DATABASE_URL=postgresql://vocalyx:vocalyx_secret@postgres:5432/vocalyx_db
Environment=REDIS_URL=redis://redis:6379/0
Environment=CELERY_BROKER_URL=redis://redis:6379/0
Environment=CELERY_RESULT_BACKEND=redis://redis:6379/0
Environment=INTERNAL_API_KEY=secret_key_pour_comms_internes_123456
Environment=ADMIN_PROJECT_NAME=ISICOMTECH
Environment=CORS_ORIGINS=http://localhost:8080,http://vocalyx-frontend:8080,http://localhost:8000,http://haproxy:8000
Environment=UPLOAD_DIR=/app/shared_uploads
Environment=LOG_LEVEL=INFO
Environment=LOG_COLORED=false
Environment=LOG_FILE_PATH=/app/logs/vocalyx-api-02.log
HealthCmd=curl -f http://localhost:8000/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=40s

[Service]
Restart=on-failure

