# vocalyx-flower.container
# ============================================================================
# VOCALYX - Flower Monitoring pour Podman/systemd
# ============================================================================
# Fichier quadlet pour le conteneur Flower (monitoring Celery)
# 
# Installation:
#   sudo cp vocalyx-flower.container /etc/containers/systemd/
#   sudo systemctl daemon-reload
#   sudo systemctl start vocalyx-flower.service
# ============================================================================

[Unit]
Description=Flower Celery Monitoring for Vocalyx
After=vocalyx-network.service vocalyx-redis.service
Requires=vocalyx-network.service vocalyx-redis.service

[Container]
Image=mher/flower:2.0
ContainerName=vocalyx-flower
RestartPolicy=unless-stopped
Network=vocalyx-network
PublishPort=5555:5555
Memory=256M
Environment=CELERY_BROKER_URL=redis://redis:6379/0
Environment=CELERY_RESULT_BACKEND=redis://redis:6379/0
Environment=FLOWER_PORT=5555
Exec=celery --broker=redis://redis:6379/0 flower --port=5555

[Service]
Restart=on-failure

