# vocalyx-redis.service
# ============================================================================
# VOCALYX - Redis Service pour systemd
# ============================================================================
# Service systemd pour gérer le conteneur Redis avec Podman
# 
# Installation:
#   sudo cp vocalyx-redis.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable vocalyx-redis.service
#   sudo systemctl start vocalyx-redis.service
# ============================================================================

[Unit]
Description=Redis Broker and Cache for Vocalyx
Documentation=https://redis.io/
After=network-online.target
Wants=network-online.target
Requires=vocalyx-redis-data.service

[Service]
Type=notify
NotifyAccess=all
TimeoutStartSec=0
Restart=always
RestartSec=10

# Créer le conteneur s'il n'existe pas, sinon le démarrer
ExecStartPre=-/usr/bin/podman rm -f vocalyx-redis
ExecStart=/usr/bin/podman run \
    --name vocalyx-redis \
    --network vocalyx-network \
    --volume vocalyx-redis-data:/data:Z \
    --publish 6379:6379 \
    --memory 1g \
    --health-cmd "redis-cli ping" \
    --health-interval 10s \
    --health-timeout 5s \
    --health-retries 5 \
    --rm \
    redis:7-alpine \
    redis-server --appendonly yes --maxmemory 256mb --maxmemory-policy allkeys-lru

ExecStop=/usr/bin/podman stop -t 10 vocalyx-redis
ExecStopPost=-/usr/bin/podman rm -f vocalyx-redis

[Install]
WantedBy=multi-user.target

