# vocalyx-enrichment-01.service
# ============================================================================
# VOCALYX - Enrichment Worker #1 Service pour systemd
# ============================================================================
# Service systemd pour gérer le conteneur vocalyx-enrichment-01 avec Podman
# 
# Note: Modifiez PROJECT_DIR dans ExecStart pour pointer vers votre projet
# 
# Installation:
#   sudo cp vocalyx-enrichment-01.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable vocalyx-enrichment-01.service
#   sudo systemctl start vocalyx-enrichment-01.service
# ============================================================================

[Unit]
Description=Vocalyx Enrichment Worker #1
After=network-online.target vocalyx-redis.service vocalyx-haproxy.service
Wants=network-online.target vocalyx-haproxy.service
Requires=vocalyx-redis.service

[Service]
Type=notify
NotifyAccess=all
TimeoutStartSec=0
Restart=always
RestartSec=10

# Modifier PROJECT_DIR pour pointer vers votre répertoire de projet
Environment="PROJECT_DIR=/home/shinohk/code/vocalyx-all"

# Créer le conteneur s'il n'existe pas, sinon le démarrer
ExecStartPre=-/usr/bin/podman rm -f vocalyx-enrichment-01
ExecStart=/usr/bin/podman run \
    --name vocalyx-enrichment-01 \
    --network vocalyx-network \
    --volume ${PROJECT_DIR}/shared/logs:/app/logs:Z \
    --volume ${PROJECT_DIR}/shared/models:/app/models:Z \
    --volume ${PROJECT_DIR}/shared/models/enrichment:/app/models/enrichment:Z \
    --cpus 4 \
    --memory 4g \
    --env INSTANCE_NAME=enrichment-worker-01 \
    --env VOCALYX_API_URL=http://haproxy:8000 \
    --env CELERY_BROKER_URL=redis://redis:6379/0 \
    --env CELERY_RESULT_BACKEND=redis://redis:6379/0 \
    --env LLM_DEVICE=cpu \
    --env INTERNAL_API_KEY=secret_key_pour_comms_internes_123456 \
    --env LOG_LEVEL=INFO \
    --env LOG_COLORED=true \
    --env LOG_FILE_PATH=/app/logs/vocalyx-enrichment-01.log \
    --rm \
    vocalyx-enrichment:latest \
    celery -A worker.celery_app worker --loglevel=info --concurrency=2 --hostname=enrichment-worker-01@%h --without-gossip --without-mingle -Q enrichment

ExecStop=/usr/bin/podman stop -t 10 vocalyx-enrichment-01
ExecStopPost=-/usr/bin/podman rm -f vocalyx-enrichment-01

[Install]
WantedBy=multi-user.target

