# vocalyx-transcribe-01.service
# ============================================================================
# VOCALYX - Transcription Worker #1 Service pour systemd
# ============================================================================
# Service systemd pour gérer le conteneur vocalyx-transcribe-01 avec Podman
# 
# Note: Modifiez PROJECT_DIR dans ExecStart pour pointer vers votre projet
# 
# Installation:
#   sudo cp vocalyx-transcribe-01.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable vocalyx-transcribe-01.service
#   sudo systemctl start vocalyx-transcribe-01.service
# ============================================================================

[Unit]
Description=Vocalyx Transcription Worker #1
After=network-online.target vocalyx-redis.service vocalyx-haproxy.service
Wants=network-online.target vocalyx-haproxy.service
Requires=vocalyx-redis.service

[Service]
Type=notify
NotifyAccess=all
TimeoutStartSec=0
Restart=always
RestartSec=10

# Modifier PROJECT_DIR pour pointer vers votre répertoire de projet
Environment="PROJECT_DIR=/home/shinohk/code/vocalyx-all"

# Créer le conteneur s'il n'existe pas, sinon le démarrer
ExecStartPre=-/usr/bin/podman rm -f vocalyx-transcribe-01
ExecStart=/usr/bin/podman run \
    --name vocalyx-transcribe-01 \
    --network vocalyx-network \
    --volume ${PROJECT_DIR}/shared/uploads:/app/shared_uploads:Z \
    --volume ${PROJECT_DIR}/shared/models:/app/models:Z \
    --volume ${PROJECT_DIR}/shared/logs:/app/logs:Z \
    --cpus 4 \
    --memory 8g \
    --env INSTANCE_NAME=worker-01 \
    --env VOCALYX_API_URL=http://haproxy:8000 \
    --env CELERY_BROKER_URL=redis://redis:6379/0 \
    --env CELERY_RESULT_BACKEND=redis://redis:6379/0 \
    --env WHISPER_MODEL=./models/openai-whisper-small \
    --env WHISPER_DEVICE=cpu \
    --env WHISPER_COMPUTE_TYPE=int8 \
    --env WHISPER_LANGUAGE=fr \
    --env MAX_WORKERS=2 \
    --env VAD_ENABLED=true \
    --env UPLOAD_DIR=/app/shared_uploads \
    --env INTERNAL_API_KEY=secret_key_pour_comms_internes_123456 \
    --env LOG_LEVEL=INFO \
    --env LOG_COLORED=true \
    --env LOG_FILE_PATH=/app/logs/vocalyx-transcribe-01.log \
    --rm \
    vocalyx-transcribe:latest \
    celery -A worker.celery_app worker --loglevel=info --concurrency=2 --hostname=worker-01@%h --without-gossip --without-mingle -Q transcription

ExecStop=/usr/bin/podman stop -t 10 vocalyx-transcribe-01
ExecStopPost=-/usr/bin/podman rm -f vocalyx-transcribe-01

[Install]
WantedBy=multi-user.target

