# vocalyx-enrichment-02.container
# ============================================================================
# VOCALYX - Enrichment Worker #2 pour Podman/systemd
# ============================================================================
# Fichier quadlet pour le conteneur vocalyx-enrichment-02
# 
# Note: L'image doit être construite au préalable:
#   podman build -t vocalyx-enrichment:latest -f ./vocalyx-enrichment/Containerfile ./vocalyx-enrichment
# 
# Installation:
#   sudo cp vocalyx-enrichment-02.container /etc/containers/systemd/
#   sudo systemctl daemon-reload
#   sudo systemctl start vocalyx-enrichment-02.service
# ============================================================================

[Unit]
Description=Vocalyx Enrichment Worker #2
After=vocalyx-network.service vocalyx-redis.service vocalyx-haproxy.service
Requires=vocalyx-network.service vocalyx-redis.service
Wants=vocalyx-haproxy.service

[Container]
Image=vocalyx-enrichment:latest
ContainerName=vocalyx-enrichment-02
RestartPolicy=unless-stopped
Network=vocalyx-network
Volume=%E/shared/logs:/app/logs:Z
Volume=%E/shared/models:/app/models:Z
Volume=%E/shared/models/enrichment:/app/models/enrichment:Z
CPUs=4
Memory=4G
Environment=INSTANCE_NAME=enrichment-worker-02
Environment=VOCALYX_API_URL=http://haproxy:8000
Environment=CELERY_BROKER_URL=redis://redis:6379/0
Environment=CELERY_RESULT_BACKEND=redis://redis:6379/0
Environment=LLM_DEVICE=cpu
Environment=INTERNAL_API_KEY=secret_key_pour_comms_internes_123456
Environment=LOG_LEVEL=INFO
Environment=LOG_COLORED=true
Environment=LOG_FILE_PATH=/app/logs/vocalyx-enrichment-02.log
Exec=celery -A worker.celery_app worker --loglevel=info --concurrency=2 --hostname=enrichment-worker-02@%%h --without-gossip --without-mingle -Q enrichment

[Service]
Restart=on-failure

